# etc/cron.hourly/enforce-share
#!/usr/bin/env bash
for user in parents famille amandine mael aedan noam teilo; do
  enforce-share /home/$user/syncthing >> %LOG_DIR%/enforce-share.log 2>&1
done
