# etc/cron.hourly/enforce-share
#!/usr/bin/env bash

source %CONF_DIR%/enforce-sharerc.sh

for share in $ENFORCED_SHARES; do
  enforce-share $share >> %LOG_DIR%/enforce-share.log 2>&1
done
